# supervisor.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=/usr/local/sbin/php-fpm

[program:worker]
command=php /path/to/your/app/artisan queue:work
numprocs=1
autostart=true
autorestart=true
